# templates/etc/profile.d/powerman.sh.erb
######### Managed by Puppet #########
# <%= @servername %>:<%= @__file__ %>
export POWERMAN_HOST=<%= @powerman_server %>
export POWERMAN_PORT=<%= @powerman_port %>

