# templates/etc/profile.d/powerman.csh.erb
######### Managed by Puppet #########

setenv POWERMAN_HOST <%= @powerman_server %>
setenv POWERMAN_PORT <%= @powerman_port %>
alias powerman "/usr/bin/powerman -h ${POWERMAN_HOST}:${POWERMAN_PORT} \!*"
alias pm "/usr/bin/pm -h ${POWERMAN_HOST}:${POWERMAN_PORT} \!*"
